# s6/alpine-root/etc/services.d/gunicorn/finish
#!/usr/bin/with-contenv bash

s6-echo "Stopping gunicorn"
killall -9 gunicorn
